# kiwi-pqc-ima/appliance/root/etc/systemd/user/test-hyperfine.timer
[Unit]
Description=Timer for hyperfine

[Timer]
OnStartupSec=20s
Unit=test-hyperfine.service

[Install]
WantedBy=timers.target
