# kiwi-pqc-ima/appliance/root/etc/systemd/user/test-hyperfine.service
[Unit]
Description=test hyperfine

[Service]
Type=exec
Restart=no
ExecStart=gnome-terminal -- /usr/local/bin/test-hyperfine.sh
